# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:nodeapp]
command=/bin/bash -c "coffee /src/server.coffee"
